# supervisord.conf
[supervisord]
nodaemon=true

[program:vncserver]
command=/usr/bin/vncserver :1 -geometry 1280x800 -depth 24
user=ubuntu
autostart=true
autorestart=true
